# deploy/career-agent.service
[Unit]
Description=Career Planning Agent - FastAPI Backend
After=network.target

[Service]
Type=simple
User=career-agent
Group=career-agent
WorkingDirectory=/opt/career-agent/app
EnvironmentFile=/opt/career-agent/app/.env
Environment=PYTHONPATH=/opt/career-agent/app
ExecStart=/opt/career-agent/venv/bin/uvicorn backend.app:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
